# src/08_generate_rootfs/etc/inittab
::sysinit:/etc/bootscript.sh
::restart:/sbin/init
::shutdown:echo "Sync file buffers..."
::shutdown:sync
::shutdown:echo "Unmount all filesystems..."
::shutdown:umount -a -r
::ctrlaltdel:/sbin/reboot
::once:cat /etc/welcome.txt
::respawn:/bin/cttyhack /bin/sh
tty2::once:cat /etc/welcome.txt
tty2::respawn:/bin/sh
tty3::once:cat /etc/welcome.txt
tty3::respawn:/bin/sh
tty4::once:cat /etc/welcome.txt
tty4::respawn:/bin/sh

